# docker/frontend/supervisord.ini
[supervisord]
nodaemon=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon

[program:nginx]
command=/usr/sbin/nginx